# alexa_bridge_admin/rootfs/etc/s6-overlay/s6-rc.d/alexa-bridge-admin/finish
#!/command/with-contenv bashio
set +e

bashio::log.warning "AlexaBridge Admin stopped"
